# config/includes.chroot/etc/skel/.config/autostart/lilydev.desktop
[Desktop Entry]
Type=Application
Name=LilyDevSetup
Exec=qterminal --workdir=~ -e ./.lilydev-setup.sh
StartupNotify=true
Terminal=true
